# chromeos-base/detachable-emulator/files/load-vdtb-module.conf
# Copyright 2019 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Load fydeos detachable device emulator before startup"
author          "yang@fydeos.io"

# This job is detecting the specific usb device to determine if is tablet mode
# enable.
#

start on stopped pre-startup

exec modprobe fyde-vdtb
